# p1-to-mqtt/docs/systemd/dsmr-mqtt.service
[Unit]
Description=DSMR to MQTT Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dsmr
SupplementaryGroups=dialout
EnvironmentFile=/etc/dsmr-mqtt.conf
ExecStart=/usr/local/bin/dsmr-mqtt
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dsmr-mqtt

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
